# lila-http.socket
[Unit]
Description=lila-http socket

[Socket]
ListenStream=/run/lila-http.socket

[Install]
WantedBy=sockets.target
